# launchagents/com.kokoro.tts-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kokoro.tts-api</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/darianrosebrook/Desktop/Projects/kokoro-onnx/.venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>api.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/darianrosebrook/Desktop/Projects/kokoro-onnx</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/darianrosebrook/Desktop/Projects/kokoro-onnx/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/darianrosebrook/Desktop/Projects/kokoro-onnx/.venv</string>
        <key>KOKORO_MODEL_FILE</key>
        <string>kokoro-v1.0.int8-graph-opt.onnx</string>
        <key>ESPEAK_DATA_PATH</key>
        <string>/Users/darianrosebrook/Desktop/Projects/kokoro-onnx/.venv/lib/python3.13/site-packages/espeakng_loader/espeak-ng-data</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/darianrosebrook/Library/Logs/kokoro/tts-api.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/darianrosebrook/Library/Logs/kokoro/tts-api.log</string>
</dict>
</plist>
